# rootfs/etc/s6-overlay/s6-rc.d/svc-mailpit/finish
#!/command/with-contenv sh
# Nothing to clean up — Mailpit handles its own shutdown
