# rootfs/etc/s6-overlay/s6-rc.d/svc-opendkim/finish
#!/command/with-contenv sh
rm -f /run/opendkim/opendkim.pid /run/opendkim/opendkim.sock 2>/dev/null || true
