# rootfs/etc/s6-overlay/s6-rc.d/svc-opendkim/run
#!/command/with-contenv sh
# Skip if DKIM is disabled or in catch mode
if [ "${DKIM_ENABLED:-false}" != "true" ] || [ "${MAIL_MODE:-catch}" = "catch" ]; then
    exec sleep infinity
fi
exec s6-setuidgid opendkim /usr/sbin/opendkim -f -l -x /etc/opendkim.conf
