# rootfs/etc/s6-overlay/s6-rc.d/svc-pilot/finish
#!/command/with-contenv sh
# Nothing to clean up
