# rootfs/etc/s6-overlay/s6-rc.d/svc-postfix/finish
#!/command/with-contenv sh
# Clean Postfix shutdown on service stop (exit code 256 = killed by signal)
[ "${MAIL_MODE:-catch}" != "catch" ] && postfix stop 2>/dev/null || true
