# rootfs/etc/s6-overlay/s6-rc.d/svc-postfix/run
#!/command/with-contenv sh
# Skip in catch mode — sleep without consuming CPU
case "${MAIL_MODE:-catch}" in
    catch)
        exec sleep infinity
        ;;
    *)
        exec postfix start-fg
        ;;
esac
